# crontab
28 1 * * * php ~/index.php file cron
28 2 * * * php ~/index.php user cron